# immich-cloud-backup.timer
[Unit]
Description=Tägliches Immich Cloud Backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target