# immich-cloud-backup.service
[Unit]
Description=Immich Cloud Backup to krDrive

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-immich.sh
